# gpu/gpu.service
[Unit]
Description=Increase GPU core and memory clocks

[Service]
Type=oneshot
ExecStart=/opt/gpu.sh
User=root

[Install]
WantedBy=multi-user.target
